# deploy/reins-deploy-agent.service
# reins deploy host agent — systemd 유닛.
# jjong(Oracle A1, user=ubuntu) 기준. 다른 환경이면 User/ExecStart 경로 수정.
#
# 설치:
#   sudo install -d /etc/reins
#   sudo cp agent.sh /etc/reins/agent.sh && sudo chmod +x /etc/reins/agent.sh   # repo 밖 안정 경로(자기 배포 중 변형 방지)
#   sudo tee /etc/reins/deploy-agent.env >/dev/null <<'EOF'
#   REINS_API_URL=http://127.0.0.1:21001
#   REINS_TOKEN=<reins API 토큰>
#   POLL_INTERVAL=5
#   EOF
#   sudo cp reins-deploy-agent.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now reins-deploy-agent
#   journalctl -u reins-deploy-agent -f   # 로그 확인
# 갱신: agent.sh 바뀌면 /etc/reins/agent.sh 재복사 후 sudo systemctl restart reins-deploy-agent.
[Unit]
Description=reins deploy host agent (board → git+deploy.sh)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/etc/reins/deploy-agent.env
ExecStart=/etc/reins/agent.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
